# templates/dev.varnish.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">

    <dict>
        <key>Label</key>
        <string>dev.varnish</string>

        <key>ProgramArguments</key>
        <array>
            <string><%= scope.lookupvar "varnish::config::executable" %></string>
            <string>-n</string>
            <string><%= scope.lookupvar "varnish::config::exec_name" %></string>
            <string>-f</string>
            <string><%= scope.lookupvar "varnish::config::vclfile" %></string>
            <string>-s</string>
            <string>malloc,1G</string>
            <string>-T</string>
            <string>127.0.0.1:2000</string>
            <string>-a</string>
            <string>0.0.0.0:<%= @port %></string>
            <string>-p</string>
            <string>pipe_timeout=3600</string>
        </array>

        <key>KeepAlive</key>
        <true/>

        <key>RunAtLoad</key>
        <true/>

        <key>WorkingDirectory</key>
        <string><%= scope.lookupvar "varnish::config::datadir" %></string>

        <key>StandardErrorPath</key>
        <string><%= scope.lookupvar "varnish::config::logfile" %></string>

        <key>StandardOutPath</key>
        <string><%= scope.lookupvar "varnish::config::logfile" %></string>

    </dict>

</plist>
